# files/system/etc/systemd/user/dotfiles.service
[Unit]
Description=Install dots

[Service]
Type=oneshot
ExecStart=/usr/bin/dots-tool
RemainAfterExit=true

[Install]
WantedBy=default.target

